# rc.local
service ntpd restart
timedatectl set-ntp true
